# supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:nodered]
command = npm start --userDir /data

[program:telldusd]
command = telldusd --nodaemon
